# templates/profile.d/pathmunge.sh.erb
# Managed by Puppet <%= @module_name %>.

# Adjust $PATH using logic similar to pathmunge() in RedHat Bourne Shell
# CSH does not support functions, so we expand this in Puppet for consistency.
<%- scope.lookupvar('etc_profile::pathmunges').each do |k,v| -%>

if [ -d <%= k %> ]; then
    case ":${PATH}:" in
        *:"<%= k %>":*)
            ;;
        *)
<%-     if v['after'] then -%>
            export PATH=$PATH:<%= k %>
<%-     else -%>
            export PATH=<%= k %>:$PATH
<%-     end -%>
            ;;
    esac
fi
<%- end -%>
