# templates/profile.d/pathmunge.csh.erb
# Managed by Puppet <%= @module_name %>.

# Adjust $PATH using logic similar to pathmunge() in RedHat Bourne Shell
# CSH does not support functions, so we expand this in Puppet for consistency.
<% scope.lookupvar('etc_profile::pathmunges').each do |k, v| %>

if ( -d <%= k %> ) then
    switch(":${PATH}:")
        case "*:<%= k %>:*":
            breaksw
        default:
<%     if v['after'] -%>
            set path = ( ${path:q} <%= k %> )
<%     else -%>
            set path = ( <%= k %> ${path:q} )
<%     end -%>
            breaksw
    endsw
endif
<% end -%>
