# templates/profile.d/usr_local.csh.erb
# Managed by Puppet <%= @module_name %>.

# Adjust $PATH using logic similar to pathmunge() in RedHat Bourne Shell
# CSH does not support functions, so we expand this in Puppet for consistency.
<%- scope.lookupvar('etc_profile::pathadd').each do |k,v| -%>

if ( -d <%= k %> ) then
    switch(":${PATH}:")
        case "*:<%= k %>:*":
            breaksw
        default:
<%-     if v['after'] then -%>
            set PATH="$PATH:<%= k %>"
<%-     else -%>
            set PATH="<%= k %>:$PATH"
<%-     end -%>
            breaksw
    endsw
endif
<%- end -%>

# Adjust $MANPATH using logic similar to pathmunge() in RedHat Bourne Shell
# CSH does not support functions, so we expand this in Puppet for consistency.
if (! $?MANPATH) then
  setenv MANPATH
endif

<%- scope.lookupvar('etc_profile::manpathadd').each do |k,v| -%>

if ( -d <%= k %> ) then
    switch(":${MANPATH}:")
        case "*:<%= k %>:*":
            breaksw
        default:
<%-     if v['after'] then -%>
            set MANPATH="$MANPATH:<%= k %>"
<%-     else -%>
            set MANPATH="<%= k %>:$MANPATH"
<%-     end -%>
            breaksw
    endsw
endif
<%- end -%>

# Ensure OS X-style command
if ( -x /usr/libexec/path_helper ) then
    eval `/usr/libexec/path_helper -c`
endif
