# templates/profile.d/usr_local.csh.erb
# Managed by Puppet.

## Add /usr/local to CSH-style environments

foreach p (<%= scope.lookupvar('etc_profile::pathadd').join(' ') %>)
    echo $PATH | fgrep $p || set PATH=${p}:${PATH}
end

foreach p (<%= scope.lookupvar('etc_profile::manpathadd').join(' ') %>)
    set MANPATH=${p}:${MANPATH}
end

# Add Puppet to the end
set PATH=${PATH}:/opt/puppetlabs/bin
set MANPATH=${MANPATH}:/opt/puppetlabs/puppet/share/man
